# debian/bitsby8.service
[Unit]
Description=BitsBy8 Serial Drive Server
Documentation=file:///usr/share/doc/bitsby8/README.md
After=network.target

[Service]
Type=simple
User=bitsby8
Group=bitsby8
SupplementaryGroups=dialout gpio

# Environment
EnvironmentFile=-/etc/default/bitsby8
WorkingDirectory=/var/lib/bitsby8

# Command
ExecStart=/usr/bin/node /usr/lib/bitsby8/dist/index.js --config /etc/bitsby8/bitsby8.config.json

# Restart policy — always, not just on-failure. The web UI's
# "Restart now" button calls process.exit(0) after a config save;
# systemd needs to bring the daemon back on a clean exit for that
# flow to work.
Restart=always
RestartSec=10s
# ...except the web UI's "Shutdown" button, which exits with code 42 to
# ask systemd to leave the daemon down (systemd-native stop, no systemctl
# call or extra privilege needed under this hardened non-root unit).
RestartPreventExitStatus=42

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# The daemon reads /etc/bitsby8/bitsby8.config.json (baseline / dpkg
# conffile) and writes UI-driven changes to
# /var/lib/bitsby8/bitsby8.overrides.json — /etc stays read-only under
# ProtectSystem=strict, as intended.
ReadWritePaths=/var/lib/bitsby8 /var/log/bitsby8
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bitsby8

[Install]
WantedBy=multi-user.target
