# debian/postinst
#!/bin/bash
set -e

# Source debconf library for maintainer scripts
. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Create fdcsds system user if it doesn't exist
        if ! getent passwd fdcsds >/dev/null; then
            adduser --system --group --home /var/lib/fdcsds \
                    --no-create-home --disabled-password \
                    --shell /usr/sbin/nologin \
                    --gecos "FDC+ Serial Drive Server" fdcsds
        fi

        # Create gpio group if it doesn't exist
        if ! getent group gpio >/dev/null; then
            addgroup --system gpio
        fi

        # Add fdcsds user to gpio group for GPIO access
        if getent group gpio >/dev/null; then
            adduser fdcsds gpio >/dev/null || true
        fi

        # Add fdcsds user to dialout group for serial port access
        if getent group dialout >/dev/null; then
            adduser fdcsds dialout >/dev/null || true
        fi

        # Set ownership and permissions on runtime directories
        chown -R fdcsds:fdcsds /var/lib/fdcsds
        chown -R fdcsds:fdcsds /var/log/fdcsds
        chmod 755 /var/lib/fdcsds
        chmod 755 /var/log/fdcsds

        # Create data subdirectories
        for dir in disks scripts uploads cassettes; do
            mkdir -p /var/lib/fdcsds/$dir
        done
        chown -R fdcsds:fdcsds /var/lib/fdcsds

        # Copy sample disk images to writable data directory (skip if already present)
        if [ -d /usr/share/fdcsds/disks ]; then
            for img in /usr/share/fdcsds/disks/*; do
                [ -f "$img" ] || continue
                base=$(basename "$img")
                if [ ! -f "/var/lib/fdcsds/disks/$base" ]; then
                    cp "$img" "/var/lib/fdcsds/disks/$base"
                    chown fdcsds:fdcsds "/var/lib/fdcsds/disks/$base"
                fi
            done
        fi

        # Create default configuration file if it doesn't exist
        if [ ! -f /etc/fdcsds/fdcsds.config.json ]; then
            cat > /etc/fdcsds/fdcsds.config.json <<'EOF'
{
  "dataDir": "/var/lib/fdcsds",
  "port": "/dev/ttyUSB0",
  "baud": 230400,
  "drive0": "disks/cpm22.dsk",
  "drive1": null,
  "drive2": null,
  "drive3": null,
  "readonly": [0],
  "verbose": false,
  "debug": false,
  "web": true,
  "webPort": 3000,
  "webHost": "0.0.0.0",
  "terminalPort": "/dev/ttyUSB1",
  "terminalBaud": 9600,
  "terminalAutoconnect": false
}
EOF
            chown root:root /etc/fdcsds/fdcsds.config.json
            chmod 644 /etc/fdcsds/fdcsds.config.json
        fi

        # Rebuild native modules against the target system's Node ABI.
        # The .deb is built once (usually on a Pi with Node 20) but pnpm's
        # store may hand out a better-sqlite3 binary compiled against a
        # different Node (Node 24 ABI 141 was the recurring offender). If
        # the shipped binary's ABI doesn't match the running Node, DB init
        # throws "compiled against a different Node.js version" and every
        # API call that touches the DB returns 500. Force-recompile in
        # place so the binary always matches whichever Node is installed.
        #
        # Runs in the pnpm virtual-store path since node-linker=hoisted
        # leaves .pnpm/<pkg>/node_modules/<pkg>/ as the source of truth
        # for the top-level symlink at node_modules/<pkg>/.
        if command -v npm >/dev/null 2>&1; then
            for pkg in better-sqlite3; do
                pkg_dir=$(find /usr/lib/fdcsds/node_modules/.pnpm -maxdepth 1 \
                    -type d -name "${pkg}@*" 2>/dev/null | head -1)
                if [ -n "$pkg_dir" ] && [ -d "$pkg_dir/node_modules/$pkg" ]; then
                    echo "Rebuilding $pkg for local Node ABI..."
                    (cd "$pkg_dir/node_modules/$pkg" && npm rebuild --build-from-source) \
                        || echo "WARNING: $pkg rebuild failed — DB features may not work."
                fi
            done
        else
            echo "WARNING: npm not on PATH — cannot rebuild native modules."
            echo "         If the DB fails to init, install npm and run:"
            echo "         cd /usr/lib/fdcsds/node_modules/.pnpm/better-sqlite3@*/node_modules/better-sqlite3 && npm rebuild"
        fi

        # Reload systemd daemon to recognize new service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null || true
        fi

        # Display post-installation message
        echo ""
        echo "========================================================================"
        echo "  FDC+ Serial Drive Server (fdcsds) has been installed successfully!"
        echo "========================================================================"
        echo ""
        echo "Next steps:"
        echo ""
        echo "1. Edit the configuration file:"
        echo "   sudo nano /etc/fdcsds/fdcsds.config.json"
        echo ""
        echo "2. Configure your serial ports and disk images"
        echo ""
        echo "3. Start the service:"
        echo "   sudo systemctl start fdcsds"
        echo ""
        echo "4. Enable automatic startup on boot:"
        echo "   sudo systemctl enable fdcsds"
        echo ""
        echo "5. Check service status:"
        echo "   sudo systemctl status fdcsds"
        echo ""
        echo "6. View logs:"
        echo "   sudo journalctl -u fdcsds -f"
        echo ""
        echo "7. Access the web interface (if enabled):"
        echo "   http://$(hostname -I | awk '{print $1}'):3000"
        echo ""
        echo "Documentation: /usr/share/doc/fdcsds/"
        echo "Sample disk images: /usr/share/fdcsds/disks/"
        echo "========================================================================"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        # Nothing to do
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
